# printer-agent/streetfood-printer-agent.service
[Unit]
Description=Street Food Events — ESC/POS Printer Agent
Documentation=https://github.com/NelloRizzo/StreetFoodEvents
After=network.target

[Service]
Type=simple
User=pi
Group=lp
WorkingDirectory=/opt/streetfood-printer-agent
ExecStart=/usr/bin/node dist/server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=9300
Environment=PRINTER_DEVICE=/dev/usb/lp0
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
