# photo-point/photo-point.service
[Unit]
Description=Street Food Events — Photo Point
Documentation=https://github.com/NelloRizzo/StreetFoodEvents
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/photo-point
ExecStart=/usr/bin/python3 -m gunicorn app:app --bind 0.0.0.0:9400 --workers 2 --timeout 30
Restart=always
RestartSec=5
Environment="PHOTO_POINT_HOST=0.0.0.0"
Environment="PHOTO_POINT_PORT=9400"
Environment="PHOTO_POINT_CAPTURES=/opt/photo-point/captures"
Environment="PHOTO_POINT_FRAMES=/opt/photo-point/frames"
Environment="PHOTO_POINT_CAMERA=0"
Environment="PHOTO_POINT_RESOLUTION=1920x1080"
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
